# xfce/autostart/plank.desktop
[Desktop Entry]
Type=Application
Name=Plank
Comment=Dock
Exec=plank
Icon=plank
X-GNOME-Autostart-enabled=true
